# nouveau-workaround/systemd/nouveau-sleep-daemon.service
[Unit]
Description=Rust daemon to inhibit sleep and toggle NVIDIA outputs
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nouveau_sleep_daemon
Restart=on-failure

[Install]
WantedBy=default.target
